# systemd/cat-tv.service
[Unit]
Description=Cat TV - YouTube Entertainment System for Cats
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/cat_tv
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="HOME=/home/pi"
ExecStart=/home/pi/.local/bin/uv run cat-tv
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give the service permission to control display
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target